# vpn-bot/vpn-bot.service
[Unit]
Description=Arengate VPN Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vpn
EnvironmentFile=-/root/vpn/.env
EnvironmentFile=-/root/vpn/.env.local
ExecStart=/usr/local/bin/vpn-bot
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
